# templates/etc/network/if-up.d/iproute2-002-rules-Debian.erb
#!/bin/sh
###
### File managed by Puppet
###

[ "$IFACE" != "lo" ] || exit 0

#
# Add rules
#
<% @rules.each do |value| -%>
/sbin/ip rule add<% if value['from'] -%> from <%= value['from'] %><% end -%><% if value['to'] -%> to <%= value['to'] %><% end -%><% if value['fwmark'] -%> fwmark <%= value['fwmark'] %><% end -%><% if value['table'] -%> table <%= value['table'] %><% end -%><% if value['priority'] -%> priority <%= value['priority'] %><% end -%>
<% end -%>
